# kali-mcp.service
[Unit]
Description=Kali MCP Server
After=network.target

[Service]
Type=simple
User=neeraj
WorkingDirectory=/home/neeraj/kali-mcp
ExecStart=/home/neeraj/kali-mcp/venv/bin/python server.py
Restart=on-failure
RestartSec=5
StandardInput=null
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
